# root/etc/s6-overlay/s6-rc.d/init-series-scraper-config/up
/etc/s6-overlay/s6-rc.d/init-series-scraper-config/run
